# bridge/remarkable/systemd/codrawer-bridge-native.service
[Unit]
Description=CoDrawer Pen Bridge (native binary, reMarkable Paper Pro)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=DESKTOP_WS=ws://192.168.50.2:8000/ws/session1
# Default to NO_GRAB=1 so local ink still works while streaming.
Environment=NO_GRAB=1
Environment=PING_SECONDS=2
Environment=PONG_TIMEOUT_SECONDS=8
ExecStart=/home/root/codrawer_bridge_native -ws ${DESKTOP_WS}
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target


